# vendors/Elster/atlas/init/etc/inittab
ttyS1::respawn:/bin/sh
null::respawn:/bin/mgate
